# upstart.conf
description "TLS policy map"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 10 5

setuid postfix

exec /usr/bin/policy-client --socket /var/spool/postfix/run/tlspolicy.sock
